# ansible/roles/tendermint-exporter/templates/tendermint_exporter_healer.service.j2
[Unit]
Description=Tendermint Exporter Healer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/tendermint_exporter_healer.sh